# s6-overlay/s6-rc.d/svc-ofelia/finish
#!/command/with-contenv bash

exec s6-setuidgid bipops /usr/bin/rsnapshot -c /etc/rsnapshot.conf latest
