# s6-rc.d/monitor-game/run
#!/command/with-contenv bash

if [ ! -f "/gameservers/$GAMESERVER/monitor.sh" ]; then
  sleep infinity
fi

exec /gameservers/$GAMESERVER/monitor.sh
