# s6-rc.d/configure-game/up
/etc/s6-overlay/s6-rc.d/configure-game/scripts/configure-game.sh